# flask-server/deploy/systemd/pants-flask.service
[Unit]
Description=PanTS Flask/Gunicorn Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=visitor
Group=visitor
WorkingDirectory=/home/visitor/PanTS-Viewer/flask-server

# Single process with threads so in-memory inference_jobs dict is shared
ExecStart=/usr/bin/gunicorn \
    --worker-class gthread \
    --workers 1 \
    --threads 8 \
    --bind 127.0.0.1:8000 \
    --timeout 3600 \
    --access-logfile - \
    --error-logfile - \
    app:app

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
